# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/SalishSea/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/SalishSea/crontab.log 2>&1
